# root/etc/services.d/irkerd/run
#!/usr/bin/with-contenv bash
set -e

: ${LISTEN_ADDR:="0.0.0.0"}
: ${LOG_LEVEL:="debug"}

exec s6-setuidgid app /usr/bin/irkerd -l /dev/stdout -H $LISTEN_ADDR -d $LOG_LEVEL
